# scripts/launchd/com.finances.watchlist-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.finances.watchlist-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/nyhzdev/devroom/battlefields/finances/scripts/cron-watchlist-sync.sh</string>
  </array>

  <!-- Intraday refresh for watchlisted assets: every 5 min. Separate lane from
       the daily 23:00 close; never writes price_history. launchd coalesces the
       ticks missed while the Mac was asleep. -->
  <key>StartInterval</key>
  <integer>300</integer>

  <!-- Run once right after load so the cache isn't stale until the first tick. -->
  <key>RunAtLoad</key>
  <true/>

  <!-- One-shot curl, not a daemon: do NOT keep alive. ThrottleInterval guards
       against a tight relaunch loop if curl fails fast. -->
  <key>ThrottleInterval</key>
  <integer>60</integer>

  <key>StandardOutPath</key>
  <string>/Users/nyhzdev/.finances/logs/watchlist-sync.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/nyhzdev/.finances/logs/watchlist-sync.log</string>
</dict>
</plist>
